# root/etc/s6-overlay/s6-rc.d/init-packt-cli-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# create folders
mkdir -p \
    /config/{backup,conf,data,data/tmp,web,logs,plugins}

# set permissions
lsiown -R abc:abc /config
