# root/etc/s6-overlay/s6-rc.d/init-packt-cli-config/up
/etc/s6-overlay/s6-rc.d/init-packt-cli-config/run
