# root/etc/crontabs/root
0 1 * * * packt-cli -gd -c /config/configFile.cfg 2>&1 | tee /config/logFile.log
