# root/var/spool/cron/crontabs/root
21 2 * * * /root/Packt-Publishing-Free-Learning/src/packtPublishingFreeEbook.py -c /config -g -da
