# root/var/spool/cron/crontabs/root
0 1 * * * packt-cli -gd -c /config/configFile.cfg > /config/logFile.log 2>&1