# root/var/spool/cron/crontabs/root
*/10 * * * * python3 /root/Packt-Publishing-Free-Learning/src/packtPublishingFreeEbook.py -c /config/configFile.cfg -gd > /config/packtPublishingFreeEbook.log
