# root/var/spool/cron/crontabs/root
35 19 * * * packt-cli -gd -c /config/configFile.cfg > /config/LOG_FILE.log 2>&1